# systemd/certbot-smart.service
[Unit]
Description=Certbot Smart Manager - SSL Certificate Renewal Service
Documentation=https://github.com/shuvo-halder/certbot-smart-manager
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ssl-manager --auto
User=root
Group=root
UMask=027
StandardOutput=append:/var/log/certbot-smart-manager.log
StandardError=append:/var/log/certbot-smart-manager.log
Environment=HOME=/root
Environment=CONFIG_FILE=/etc/certbot-smart-manager/certbot-smart.conf

# Security hardening
LockPersonality=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
ReadWritePaths=/etc/letsencrypt /var/log/certbot-smart-manager.log /var/lock /var/backups/certbot-smart-manager
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW CAP_DAC_OVERRIDE CAP_CHOWN CAP_FOWNER CAP_SETUID CAP_SETGID
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target